# crontab
MAILTO=""
52 * * * * jsub -cwd -N sdbot run.job
